# service/survey.service
[Unit]
Description=Survey host

[Service]
Type=forking
WorkingDir=/opt/epiphyte/survey/
ExecStart=/usr/bin/python survey.py
