# service/survey.service
[Unit]
Description=Survey host

[Service]
WorkingDir=/opt/epiphyte/survey
Type=forking
ExecStart=/usr/bin/survey
