# service/survey.service
[Unit]
Description=Survey host

[Service]
Type=forking
ExecStart=/usr/bin/survey
