# service/survey.service
[Unit]
Description=Survey host

[Service]
Type=simple
ExecStart=/usr/bin/survey
